# thakhin-os/config/includes.chroot/etc/profile.d/99-thakhin.sh
# Thakhin OS shell conveniences
alias facts='thakhin-facts'
alias sayar-gyi='thakhin-facts'
alias quiz='tarzan-check'
alias ukkahta='xdg-open /usr/share/thakhin-os/website/index.html 2>/dev/null || echo "See: /usr/share/thakhin-os/website/index.html"'
